# files/supervisord-salt.conf
[supervisord]
nodaemon=true


[program:salt-master]
command=/usr/bin/salt-master
autorestart=true

